# usb-backup.service
[Unit]
Description=USB Backup Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb_backup.py
