# usb-backup.timer
[Unit]
Description=Run USB Backup Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
